# lib/windows/queries/windows_anomaly_process-execution.yml
- name: InstallUtil Execute, InstallUtil is a command-line utility that allows for installation and uninstallation of resources by executing specific installer components specified in .NET binaries
  automations_enabled: false
  description: InstallUtil Execute, InstallUtil is a command-line utility that allows for installation and uninstallation of resources by executing specific installer components specified in .NET binaries - ATT&CK T1118
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%InstallUtil%';
- name: PsExec Execute, is a free Microsoft tool that can be used to execute a program on another computer.
  automations_enabled: false
  description: PsExec Execute, is a free Microsoft tool that can be used to execute a program on another computer. - ATT&CK T1035,S0029
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%PsExec%';
- name: Snapshot Prefetach File Directory
  automations_enabled: false
  description: Snapshot Prefetach File Directory - ATT&CK T1107
  discard_data: false
  interval: 28800
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\';
- name: Schtasks Execute, usaullay used to create a scheduled task
  automations_enabled: false
  description: Schtasks Execute, usaullay used to create a scheduled task - ATT&CK T1053,S0110
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%at.exe%';
- name: Attrib Execute, usaullay used to modify file attributes
  automations_enabled: false
  description: Attrib Execute, usaullay used to modify file attributes - ATT&CK T1158
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%attrib%';
- name: Bitsadmin Execute, Windows Background Intelligent Transfer Service (BITS) is a low-bandwidth, asynchronous file transfer mechanism exposed through Component Object Model (COM)
  automations_enabled: false
  description: Bitsadmin Execute, Windows Background Intelligent Transfer Service (BITS) is a low-bandwidth, asynchronous file transfer mechanism exposed through Component Object Model (COM) - ATT&CK T1197,S0190
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%bitsadmin%';
- name: Certutil Execute, Certutil.exe is a legitimate built-in command-line program to manage certificates in Windows
  automations_enabled: false
  description: Certutil Execute, Certutil.exe is a legitimate built-in command-line program to manage certificates in Windows - ATT&CK T1105,T1140,T1130,S0160
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%certutil%';
- name: Command-Line Interface Execute, CMD execution
  automations_enabled: false
  description: Command-Line Interface Execute, CMD execution - ATT&CK T1059
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%CMD%';
- name: CMSTP Execute, The Microsoft Connection Manager Profile Installer (CMSTP.exe) is a command-line program used to install Connection Manager service profiles.
  automations_enabled: false
  description: CMSTP Execute, The Microsoft Connection Manager Profile Installer (CMSTP.exe) is a command-line program used to install Connection Manager service profiles. - ATT&CK T1191
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%CMSTP%';
- name: Command-Line Interface Execute, Cscript execution starts a script so that it runs in a command-line environment.
  automations_enabled: false
  description: Command-Line Interface Execute, Cscript execution starts a script so that it runs in a command-line environment. - ATT&CK T1216
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%cscript%';
- name: Esentutl Execute, is a legitimate built-in command-line program it could be used to create a exe from dump raw source.
  automations_enabled: false
  description: Esentutl Execute, is a legitimate built-in command-line program it could be used to create a exe from dump raw source.
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%esentutl%';
- name: Mshta Execute, is a utility that executes Microsoft HTML Applications (HTA)
  automations_enabled: false
  description: Mshta Execute, is a utility that executes Microsoft HTML Applications (HTA) - ATT&CK T1170
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%mshta%';
- name: mstsc.exe Execute, usaullay used to perform a RDP Session 
  automations_enabled: false
  description: mstsc.exe Execute, usaullay used to perform a RDP Session  - ATT&CK T1076
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%mstsc%';
- name: Net Execute, is used in command-line operations for control of users, groups, services, and network connections
  automations_enabled: false
  description: Net Execute, is used in command-line operations for control of users, groups, services, and network connections - ATT&CK T1126,T1087,T1201,T1069,S0039,T1018,T1007,T1124
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%net.exe%';
- name: Netsh Execute, Netsh.exe (also referred to as Netshell) is a command-line scripting utility used to interact with the network configuration of a system
  automations_enabled: false
  description: Netsh Execute, Netsh.exe (also referred to as Netshell) is a command-line scripting utility used to interact with the network configuration of a system - ATT&CK T1128,T1063,S0108
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%netsh%';
- name: Netstat Execute,  is an operating system utility that displays active TCP connections, listening ports, and network statistics.
  automations_enabled: false
  description: Netstat Execute,  is an operating system utility that displays active TCP connections, listening ports, and network statistics. - ATT&CK T1049,S0104
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%netstat%';
- name: POWERSHELL Execute, is a powerful interactive command-line interface and scripting environment included in the Windows operating system
  automations_enabled: false
  description: POWERSHELL Execute, is a powerful interactive command-line interface and scripting environment included in the Windows operating system - ATT&CK T1086
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%POWERSHELL%';
- name: Reg Execute, Reg is a Windows utility used to interact with the Windows Registry. 
  automations_enabled: false
  description: Reg Execute, Reg is a Windows utility used to interact with the Windows Registry.  - ATT&CK T1214,T1012,T1063,S0075
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%reg%';
- name: Regedit Execute, is a Windows utility used to interact with the Windows Registry.
  automations_enabled: false
  description: Regedit Execute, is a Windows utility used to interact with the Windows Registry. - ATT&CK T1214
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%regedit%';
- name: POWERSHELL Execute, is a powerful interactive command-line interface and scripting environment included in the Windows operating system - Regsvr32 Detection
  automations_enabled: false
  description: POWERSHELL Execute, is a powerful interactive command-line interface and scripting environment included in the Windows operating system - ATT&CK T1117
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%regsvr%';
- name: Runas Execute, Allows a user to run specific tools and programs with different permissions than the user's current logon provides.
  automations_enabled: false
  description: Runas Execute, Allows a user to run specific tools and programs with different permissions than the user's current logon provides. - ATT&CK T1134
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%runas%';
- name: SC.exe Execute, Service Control - Create, Start, Stop, Query or Delete any Windows SERVICE. .
  automations_enabled: false
  description: SC.exe Execute, Service Control - Create, Start, Stop, Query or Delete any Windows SERVICE. . - ATT&CK T1007
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%sc%';
- name: Schtasks Execute, usaullay used to create a scheduled task - Secondary Detection
  automations_enabled: false
  description: Schtasks Execute, usaullay used to create a scheduled task - ATT&CK T1053,S0111
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%schtasks%';
- name: SVCHOST Processes not using the -k [name] convention
  automations_enabled: false
  description: SVCHOST Processes not using the -k [name] convention
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from processes where name like 'svchost.exe' and cmdline not like '%-k%';
- name: Systeminfo Execute, Systeminfo is a Windows utility that can be used to gather detailed information about a computer.
  automations_enabled: false
  description: Systeminfo Execute, Systeminfo is a Windows utility that can be used to gather detailed information about a computer. - ATT&CK T1082,S0096
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%systeminfo%';
- name: taskeng Execute, usaullay used to create a scheduled task
  automations_enabled: false
  description: taskeng Execute, usaullay used to create a scheduled task - ATT&CK T1053
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%taskeng%';
- name: Taskkill Execute, usaullay used to kill task 
  automations_enabled: false
  description: Taskkill Execute, usaullay used to kill task 
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%taskkill%';
- name: Tasklist Execute, usaullay used to list task
  automations_enabled: false
  description: Tasklist Execute, usaullay used to list task - ATT&CK T1057,T1063,T1007,S0057
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%tasklist%';
- name: tscon.exe Execute, usaullay used to Terminal Services Console
  automations_enabled: false
  description: tscon.exe Execute, usaullay used to Terminal Services Console - ATT&CK T1076
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%tscon%';
- name: Vssadmin Execute, usaullay used to execute activity on Volume Shadow copy
  automations_enabled: false
  description: Vssadmin Execute, usaullay used to execute activity on Volume Shadow copy
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%vssadmin%';
- name: Whoami Execute, used to prints the effective username of the current user
  automations_enabled: false
  description: Whoami Execute, used to prints the effective username of the current user
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%whoami%';
- name: Xcopy Execute, is used for copying multiple files or entire directory trees from one directory to another and for copying files across a network.
  automations_enabled: false
  description: Xcopy Execute, is used for copying multiple files or entire directory trees from one directory to another and for copying files across a network.
  discard_data: false
  interval: 600
  logging: 
  observer_can_run: false
  platform: windows
  query: select * from file WHERE directory = 'C:\Windows\Prefetch\' and filename like '%xcopy%';
