# lib/macos/policies/CIS_-_Ensure_XProtect_Is_Running_and_Updated.yml
- name: CIS - Ensure XProtect Is Running and Updated
  query: |
    SELECT 1
    WHERE (
        SELECT COUNT(*)
        FROM launchd
        WHERE path IN (
            '/Library/Apple/System/Library/LaunchDaemons/com.apple.XprotectFramework.PluginService.plist',
            '/Library/Apple/System/Library/LaunchDaemons/com.apple.XProtect.daemon.scan.plist'
        )
    ) = 2;
  critical: false
  description: |
    XProtect is Apple's native signature-based antivirus technology. XProtect both finds and blocks the execution of known malware. There are many AV and Endpoint Threat Detection and Response (ETDR) tools available for Mac OS. The native Apple provisioned tool looks for specific known malware and is completely integrated into the OS. No matter what other tools are being used, XProtect should have the latest signatures available.
  resolution: |
    Ask your system administrator to deploy a script that will configure:
    /usr/bin/sudo /bin/launchctl load -w /Library/Apple/System/Library/LaunchDaemons/com.apple.XProtect.daemon.scan.plist
    /usr/bin/sudo /bin/launchctl load -w /Library/Apple/System/Library/LaunchDaemons/com.apple.XprotectFramework.PluginService.plist
  platform: darwin
